# old_files_I_need_to_go_through/.config/systemd/user/default.target.wants/downloader.service
[Unit]
Description=musicazoo downloader service

[Service]
EnvironmentFile=/home/musicazoo/musicazooing/config.env
ExecStart=/usr/bin/python3 /home/musicazoo/musicazooing/downloader.py
Restart=always

[Install]
WantedBy=default.target
