# old_files_I_need_to_go_through/.config/systemd/user/wipe.service
[Unit]
Description=musicazoo wipe service

[Service]
ExecStart=/home/musicazoo/musicazooing/wipe.sh
Restart=never

[Install]
WantedBy=default.target
